# systemd/adhdaemon.socket
[Socket]
ListenDatagram = /tmp/adh-rs.socket
FileDescriptorName = noise
Service = adhdaemon.service

[Install]
WantedBy = sockets.target
